# scripts/thermalright-lcd-control.service
[Unit]
Description=Thermalright LCD Control Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=@user_home@/.local/share/thermalright-lcd-control
Environment=PYTHONPATH=@user_home@/.local/share/thermalright-lcd-control
ExecStart=thermalright-lcd-control-service
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
